# redis/templates/redis.service.j2
[Unit]
Description=Reddis In-Memory Data Store
After=network.target

[Service]
Type=forking
ExecStart={{ redis_root_path }}/bin/redis-server 
